# service/nginx-log-forwarder/run
#!/bin/bash

exec  setuidgid www-data tail -F -n 0 /var/log/nginx/access.log