# service/nginx/run
#!/bin/bash

exec  /usr/sbin/nginx -c /etc/service/nginx/nginx.conf  -g "daemon off;"