# service/nginx/run
#!/bin/bash

exec  setuidgid www-data /usr/sbin/nginx -c /etc/service/nginx/nginx.conf  -g "daemon off;"