# service/nginx/run
#!/bin/bash

sv -w3 check jussi
exec setuidgid www-data /usr/sbin/nginx -c /etc/service/nginx/nginx.conf  -g "daemon off;"