# service/scalyr-agent-2/run
#!/bin/bash -e

CONFIG=/etc/service/scalyr-agent-2/agent.json

if [[ -n ${SCALYR_API_KEY+x} ]]; then
    exec 2>&1 scalyr-agent-2 --no-fork --config-file="${CONFIG}" start
fi
