# service/jussi/run
#!/bin/bash

exec setuidgid nobody python3 "${WSGI_APP}"