# service/jussi/run
#!/bin/bash -e

cd "${APP_ROOT}"
exec 2>&1 \
  chpst -u www-data \
    pipenv run python3 "${WSGI_APP}" \
      --server_path "${HTTP_UPSTREAM_PYTHON_SERVER_PATH}" 2>&1